# base-nginx-php5-mongodb-mysql/supervisord-mongod.conf
[program:mongod]
command=/usr/bin/mongod --setParameter failIndexKeyTooLong=false
autorestart=true