# base/supervisord-php.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm --nodaemonize
username = www-data
autostart = false
autorestart = true