# ontap/supervisord-beanstalkd.conf
[program:beanstalkd]
command=/usr/bin/beanstalkd
autorestart=true